# templates/upstart.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

description "start/stop/control collectd"
start on started networking and filesystem
stop on runlevel [!2345]

env DAEMON={{collectd_prefix}}/sbin/collectd

# Tell upstart to watch for forking when tracking the pid for us.
expect fork

# prevent thrashing - 10 restarts in 5 seconds
respawn
respawn limit 10 5

# Make a log available in /var/log/upstart/collectd.log
console log

# The daemon will not start if the configuration is invalid.
pre-start exec $DAEMON -t
# Let's Fork!
exec $DAEMON
